# deploy/cloud-sql-proxy.service
[Unit]
Description=Cloud SQL Auth Proxy
After=network.target

[Service]
Type=simple
DynamicUser=yes
ExecStart=/usr/local/bin/cloud-sql-proxy rustytalon-prod:us-central1:rustytalon-db --port=5432
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
