# deploy/rustytalon.service
[Unit]
Description=RustyTalon AI Assistant
After=cloud-sql-proxy.service docker.service
Requires=cloud-sql-proxy.service

[Service]
Type=simple
ExecStartPre=/usr/bin/docker pull us-central1-docker.pkg.dev/rustytalon-prod/rustytalon/agent:latest
ExecStart=/usr/bin/docker run --rm \
  --name rustytalon \
  --env-file /opt/rustytalon/.env \
  --network=host \
  us-central1-docker.pkg.dev/rustytalon-prod/rustytalon/agent:latest
ExecStop=/usr/bin/docker stop rustytalon
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
